# nursery/obsidian-user-base/launchd/com.joeblack.wiki-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.joeblack.wiki-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <!-- uv run resolves the venv automatically -->
        <string>/Users/joeblack/.local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>src/watcher.py</string>
    </array>

    <!-- Run from the repo so CLAUDE.md is auto-loaded by claude CLI -->
    <key>WorkingDirectory</key>
    <string>/Users/joeblack/github/joeblackwaslike/idea-nursery/projects/obsidian-user-base</string>

    <!-- Start at login, restart on crash -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- stdout/stderr go to rotating logs -->
    <key>StandardOutPath</key>
    <string>/Users/joeblack/Library/Logs/wiki-watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/joeblack/Library/Logs/wiki-watcher-error.log</string>

    <!-- Inherit PATH so claude CLI is found -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/joeblack/.nvm/versions/node/v24.10.0/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/joeblack</string>
    </dict>
</dict>
</plist>
