# projects/openclaw-setup/launchd/ai.openclaw.env-setup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ai.openclaw.env-setup.plist
  Sets CLAUDE_OAUTH_TOKEN in launchd global env by reading from .secrets.env.
  install.sh substitutes __USER__
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.env-setup</string>
    <key>Comment</key>
    <string>Sets CLAUDE_OAUTH_TOKEN in launchd global env from .secrets.env</string>
    <key>RunAtLoad</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/__USER__/.openclaw/scripts/env-setup.sh</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/__USER__/.openclaw/logs/env-setup.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/__USER__/.openclaw/logs/env-setup.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/__USER__</string>
    </dict>
  </dict>
</plist>
