# nursery/openclaw-setup/launchd/ai.openclaw.secrets-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- install.sh substitutes __USER__, __OP_SA_TOKEN__ -->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.secrets-refresh</string>
    <key>Comment</key>
    <string>Refresh OpenClaw secrets from 1Password daily at 9am</string>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/__USER__/.openclaw/refresh-secrets.sh</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/__USER__/.openclaw/logs/secrets-refresh.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/__USER__/.openclaw/logs/secrets-refresh.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/__USER__</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
      <key>OP_SERVICE_ACCOUNT_TOKEN</key>
      <string>__OP_SA_TOKEN__</string>
      <key>OP_BIOMETRIC_UNLOCK_ENABLED</key>
      <string>false</string>
    </dict>
  </dict>
</plist>
