# nursery/openclaw-setup/launchd/ai.openclaw.secrets-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- install.sh substitutes __USER__, __OP_SA_TOKEN__ -->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.secrets-watchdog</string>
    <key>Comment</key>
    <string>Watch gateway error log; refresh secrets on auth failure</string>
    <key>WatchPaths</key>
    <array>
      <string>/Users/__USER__/.openclaw/logs/gateway.err.log</string>
    </array>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/__USER__/.openclaw/scripts/secrets-watchdog.sh</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/__USER__/.openclaw/logs/secrets-watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/__USER__/.openclaw/logs/secrets-watchdog.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/__USER__</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
      <key>OP_SERVICE_ACCOUNT_TOKEN</key>
      <string>__OP_SA_TOKEN__</string>
      <key>OP_BIOMETRIC_UNLOCK_ENABLED</key>
      <string>false</string>
    </dict>
  </dict>
</plist>
