# nursery/openclaw-setup/launchd/ai.openclaw.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  ai.openclaw.gateway.plist
  Installed to: ~/Library/LaunchAgents/
  install.sh substitutes __USER__, __TMPDIR__, __OP_SA_TOKEN__, __OPENCLAW_VERSION__
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.gateway</string>

    <key>Comment</key>
    <string>OpenClaw Gateway (v__OPENCLAW_VERSION__)</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>Umask</key>
    <integer>63</integer>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/__USER__/.openclaw/start-gateway.sh</string>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/__USER__/.openclaw/logs/gateway.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/__USER__/.openclaw/logs/gateway.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/__USER__</string>
      <key>TMPDIR</key>
      <string>__TMPDIR__</string>
      <key>NODE_EXTRA_CA_CERTS</key>
      <string>/etc/ssl/cert.pem</string>
      <key>NODE_USE_SYSTEM_CA</key>
      <string>1</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/Users/__USER__/Library/pnpm:/Users/__USER__/.local/bin:/Users/__USER__/bin:/usr/local/bin:/usr/bin:/bin</string>
      <key>OPENCLAW_GATEWAY_PORT</key>
      <string>18789</string>
      <key>OPENCLAW_LAUNCHD_LABEL</key>
      <string>ai.openclaw.gateway</string>
      <key>OPENCLAW_SYSTEMD_UNIT</key>
      <string>openclaw-gateway.service</string>
      <key>OPENCLAW_WINDOWS_TASK_NAME</key>
      <string>OpenClaw Gateway</string>
      <key>OPENCLAW_SERVICE_MARKER</key>
      <string>openclaw</string>
      <key>OPENCLAW_SERVICE_KIND</key>
      <string>gateway</string>
      <key>OPENCLAW_SERVICE_VERSION</key>
      <string>__OPENCLAW_VERSION__</string>
      <key>OP_SERVICE_ACCOUNT_TOKEN</key>
      <string>__OP_SA_TOKEN__</string>
      <key>OP_BIOMETRIC_UNLOCK_ENABLED</key>
      <string>false</string>
    </dict>
  </dict>
</plist>
